# homeglo/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Home Assistant WebSocket listener
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# ---------------------------------------------------------------
# 1. Detect mode
#    - If SUPERVISOR_TOKEN is present → ADDON
#    - Else → STANDALONE
# ---------------------------------------------------------------
if [ -n "${SUPERVISOR_TOKEN:-}" ]; then
    MODE="ADDON"
    bashio::log.info "Running in addon mode (SUPERVISOR_TOKEN present)"
    
    # In addon mode, use supervisor proxy URLs
    export HA_WEBSOCKET_URL="ws://supervisor/core/api/websocket"
    export HA_REST_URL="http://supervisor/core/api"
    export HA_TOKEN="${SUPERVISOR_TOKEN}"
    
    # These are not used in URL mode but kept for compatibility
    export HA_HOST="supervisor"
    export HA_PORT="80"
    export HA_USE_SSL="false"
else
    MODE="STANDALONE"
    bashio::log.info "Running in standalone mode (no supervisor token)"
    
    # Try to get config from bashio first, then fall back to env vars
    CONFIG_HOST=$(bashio::config 'host')
    CONFIG_PORT=$(bashio::config 'port')
    CONFIG_TOKEN=$(bashio::config 'token')
    CONFIG_USE_SSL=$(bashio::config 'use_ssl')
    
    # Set host
    if [ "${CONFIG_HOST}" != "null" ] && [ -n "${CONFIG_HOST}" ]; then
        export HA_HOST="${CONFIG_HOST}"
    else
        export HA_HOST="${HA_HOST:-localhost}"
    fi
    
    # Set port
    if [ "${CONFIG_PORT}" != "null" ] && [ -n "${CONFIG_PORT}" ]; then
        export HA_PORT="${CONFIG_PORT}"
    else
        export HA_PORT="${HA_PORT:-8123}"
    fi
    
    # Set token
    if [ "${CONFIG_TOKEN}" != "null" ] && [ -n "${CONFIG_TOKEN}" ]; then
        export HA_TOKEN="${CONFIG_TOKEN}"
    else
        export HA_TOKEN="${HA_TOKEN:-}"
    fi
    
    # Set SSL
    if [ "${CONFIG_USE_SSL}" != "null" ] && [ -n "${CONFIG_USE_SSL}" ]; then
        export HA_USE_SSL="${CONFIG_USE_SSL}"
    else
        export HA_USE_SSL="${HA_USE_SSL:-false}"
    fi
    
    # Build WebSocket URL
    if [ "${HA_USE_SSL}" = "true" ]; then
        export HA_WEBSOCKET_URL="wss://${HA_HOST}:${HA_PORT}/api/websocket"
    else
        export HA_WEBSOCKET_URL="ws://${HA_HOST}:${HA_PORT}/api/websocket"
    fi
fi

# Check if token is set
if [ -z "${HA_TOKEN}" ]; then
    bashio::log.error "No Home Assistant token provided!"
    bashio::log.info "Set token in addon configuration or HA_TOKEN environment variable"
    exit 1
fi

# Log startup
bashio::log.info "Starting Home Assistant WebSocket listener..."
bashio::log.info "Mode: ${MODE}"
bashio::log.info "WebSocket URL: ${HA_WEBSOCKET_URL}"

# Get color_mode configuration
COLOR_MODE=$(bashio::config 'color_mode')
if [ "${COLOR_MODE}" != "null" ] && [ -n "${COLOR_MODE}" ]; then
    export COLOR_MODE="${COLOR_MODE}"
else
    export COLOR_MODE="rgb"
fi

# Get min_color_temp configuration
MIN_COLOR_TEMP=$(bashio::config 'min_color_temp')
if [ "${MIN_COLOR_TEMP}" != "null" ] && [ -n "${MIN_COLOR_TEMP}" ]; then
    export MIN_COLOR_TEMP="${MIN_COLOR_TEMP}"
else
    export MIN_COLOR_TEMP="500"
fi

# Get max_color_temp configuration
MAX_COLOR_TEMP=$(bashio::config 'max_color_temp')
if [ "${MAX_COLOR_TEMP}" != "null" ] && [ -n "${MAX_COLOR_TEMP}" ]; then
    export MAX_COLOR_TEMP="${MAX_COLOR_TEMP}"
else
    export MAX_COLOR_TEMP="6500"
fi

bashio::log.info "Color mode: ${COLOR_MODE}"
bashio::log.info "Color temperature range: ${MIN_COLOR_TEMP}K - ${MAX_COLOR_TEMP}K"

# Run the Python app
cd /app
exec python3 main.py
